# content/etc/profile.d/rundeck.sh
test -r /opt/rundeck/etc/profile && source /opt/rundeck/etc/profile >/dev/null
